# services/nginx/run
#! /usr/bin/env bash
exec $REPL_HOME/nginx/sbin/nginx -p $REPL_HOME/nginx