# services/nodejs/run
#! /bin/sh
cd $REPL_HOME
export listen_host="127.0.0.1"
exec nodemon --ignore public/ -e js --exec "node index.js"